# file-viewer.service
[Unit]
Description=File Viewer Backend Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/file-viewer
ExecStart=/usr/bin/python3 /home/file-viewer/file-viewer-server.py
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target